# packaging/systemd/cloudflare-dns-syncd-updater.service
[Unit]
Description=Cloudflare DNS Sync — auto-updater (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/cpanel/3rdparty/cloudflare-dns-sync/bin/cfsync update
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
