# packaging/systemd/cloudflare-dns-syncd-updater.timer
[Unit]
Description=Cloudflare DNS Sync — daily auto-update check
Documentation=https://github.com/BusiRocket/cpanel-cloudflare-dns-sync

[Timer]
OnCalendar=daily
RandomizedDelaySec=3h
Persistent=true
Unit=cloudflare-dns-syncd-updater.service

[Install]
WantedBy=timers.target
